# host/systemd/marnarmon-collector.service
[Unit]
Description=MarNarMon metrics collector (one cycle)
Documentation=https://github.com/heldmar/marnarmon
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
Group=__USER__
Environment=MARNARMON_CONFIG=/etc/marnarmon/config.yml
ExecStart=__PREFIX__/venv/bin/python -m marnarmon.collect
WorkingDirectory=__PREFIX__
# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=__DB_DIR__
