# host/systemd/marnarmon-api.service
[Unit]
Description=MarNarMon metrics API (FastAPI/uvicorn)
Documentation=https://github.com/heldmar/marnarmon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
Environment=MARNARMON_CONFIG=/etc/marnarmon/config.yml
WorkingDirectory=__PREFIX__
ExecStart=__PREFIX__/venv/bin/uvicorn marnarmon.api:app --host __API_HOST__ --port __API_PORT__
Restart=on-failure
RestartSec=5
# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=__DB_DIR__
# Server Logs opt-in: install.sh replaces the next line with
# 'SupplementaryGroups=systemd-journal' when logs are enabled, or removes it.
#__LOG_GROUP_LINE__#
# Docker Monitor opt-in: install.sh replaces the next line with
# 'SupplementaryGroups=docker' when Docker Monitor is enabled, or removes it.
# systemd merges repeated SupplementaryGroups= directives, so this composes
# additively with the journal line above when both features are enabled.
# WARNING: 'docker' group membership is effectively root-equivalent.
#__DOCKER_GROUP_LINE__#

[Install]
WantedBy=multi-user.target
