# host/systemd/marnarmon-collector.timer
[Unit]
Description=Run MarNarMon collector every __INTERVAL__ minutes
Documentation=https://github.com/heldmar/marnarmon

[Timer]
# First run shortly after boot, then on a fixed cadence.
OnBootSec=1min
OnUnitActiveSec=__INTERVAL__min
AccuracySec=15s
Persistent=true
Unit=marnarmon-collector.service

[Install]
WantedBy=timers.target
